# configs/llaman-backup.service
[Unit]
Description="LLaMan Backup"

[Service]
Type=oneshot
ExecStart=/usr/bin/llaman -ba

